# rootfs/etc/services.d/cron/run
#!/command/execlineb -P
with-contenv
exec crond -f -L /logs/cron.log
