# rootfs/etc/s6-overlay/s6-rc.d/beets-web/run
#!/command/with-contenv bash
exec s6-setuidgid 99:100 /opt/beets/bin/beet web -c "${BEETS_CONFIG}"
