# rootfs/etc/s6-overlay/s6-rc.d/beets-web/run
#!/command/with-contenv sh
exec 2>&1
exec beet web -c "$BEETS_CONFIG"
