# rootfs/etc/s6-overlay/s6-rc.d/beets-web/up
fdmove -c 2 1 s6-setuidgid 99:100 /etc/s6-overlay/s6-rc.d/beets-web/run
