# rootfs/etc/services.d/beets-web/run
#!/command/execlineb -P
with-contenv
exec beet web -c "$BEETS_CONFIG"
